# homeassistant/rootfs/etc/services.d/galaxy2mqtt/run
#!/usr/bin/env bash

exec /share/galaxy2mqtt/run
